# test/supervisord.conf
[supervisord]
nodaemon=true
redirect_stderr=true
stdout_logfile=/proc/1/fd/1
logfile_maxbytes=50MB
logfile_backups=10
umask=022
user=root
[program:Procfile_4]
command=mailrig --task tasks/move1.json
stdout_logfile=Procfile_4.log
